# Cron/crontab
* * * * * /usr/local/bin/php /var/www/html/admin/cli/cron.php > /var/log/cron/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job
